# scripts/coreos-systemd/watcher.timer
[Unit]
Description=Run watcher every minute

[Timer]
OnCalendar=*:0/1
