# scripts/coreos-systemd/watcher.service
[Unit]
Description=Prints date into /tmp/date file

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c '/home/core/check.sh'
